# gzhu/etc/hotplug.d/iface/98-minieap
#!/bin/sh
[ "$ACTION" = ifup ] || exit 0
#禁用手动启用minieap
#minieap